# config/sonar_poller.conf
[program:sonar_poller]
directory=/usr/share/sonar_poller/src
command=/usr/bin/php /usr/share/sonar_poller/src/poller.php
priority=1
autostart=true
autorestart=true
